# board/batocera/allwinner/h700/rg35xx-plus/fsoverlay/etc/udev/rules.d/60-hdmi.rules
SUBSYSTEM=="extcon", ACTION=="add|change", ENV{DEVPATH}=="/devices/platform/soc/6000000.hdmi/extcon/hdmi", RUN+="/usr/bin/rg35xx-plus-switch-output"